# bbb-buildroot/gogadgets.ini
[program:gadgets]
command=/usr/local/bin/gogadgets -g /etc/gogadgets/config.json